# laravel-worker.conf
[program:websockets]
command=/usr/bin/php /home/laravel-echo/virtual-chatpanel.herokuapp.com/artisan websockets:serve --port=9001
numprocs=1
autostart=true
autorestart=true
user=laravel-echo