# laravel-worker.conf
[program:websockets]
command=/usr/bin/php /home/laravel-echo/virtual-chatpanel.herokuapp.com/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=laravel-echo